# createlink.conf
[program:createlink]
command=/createlink.py
autostart=true
autorestart=true
startretries=1
